# com.moebot.trendspider.ticker.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>EnvironmentVariables</key>
        <dict>
                <key>PATH</key>
                <string>/Users/stephenbae/Projects/moe-bot/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
                <key>PYTHONPATH</key>
                <string>/Users/stephenbae/Projects/moe-bot</string>
        </dict>
        <key>ExitTimeOut</key>
        <integer>7200</integer>
        <key>KeepAlive</key>
        <false/>
        <key>Label</key>
        <string>com.stephenbae.moebot.trendspider.ticker.monitor</string>
        <key>ProcessType</key>
        <string>Background</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/bin/caffeinate</string>
                <string>-i</string>
                <string>/Users/stephenbae/Projects/moe-bot/plist_scripts/launchd_error_wrapper.sh</string>
                <string>trendspider-ticker-monitor</string>
                <string>/Users/stephenbae/Projects/moe-bot/.venv/bin/python</string>
                <string>/Users/stephenbae/Projects/moe-bot-trendspider-data/monitor_ticker_changes.py</string>
                <string>--full-refresh</string>
        </array>
        <key>RunAtLoad</key>
        <false/>
        <key>StandardErrorPath</key>
        <string>/Users/stephenbae/Projects/moe-bot/logs/ticker_monitor_stderr.log</string>
        <key>StandardOutPath</key>
        <string>/Users/stephenbae/Projects/moe-bot/logs/ticker_monitor_stdout.log</string>
        <key>StartCalendarInterval</key>
        <array>
                <dict>
                        <key>Hour</key>
                        <integer>20</integer>
                        <key>Minute</key>
                        <integer>30</integer>
                        <key>Weekday</key>
                        <integer>1</integer>
                </dict>
                <dict>
                        <key>Hour</key>
                        <integer>20</integer>
                        <key>Minute</key>
                        <integer>30</integer>
                        <key>Weekday</key>
                        <integer>3</integer>
                </dict>
                <dict>
                        <key>Hour</key>
                        <integer>20</integer>
                        <key>Minute</key>
                        <integer>30</integer>
                        <key>Weekday</key>
                        <integer>5</integer>
                </dict>
        </array>
        <key>WorkingDirectory</key>
        <string>/Users/stephenbae/Projects/moe-bot-trendspider-data</string>
</dict>
</plist>
